# flask_setup/starter/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=gunicorn -b 0.0.0.0:80 main:app
stdout_logfile=/flask.log